# ipsec-xfrm
#!/sbin/openrc-run

extra_commands="flush"
extra_started_commands="update"

IPSEC_HOSTS_STABLE=${IPSEC_HOSTS_STABLE:-/var/lib/ipsec/ipsechosts}
IPSEC_HOSTS_UPDATE=${IPSEC_HOSTS_UPDATE:-/var/lib/ipsec/ipsechosts-update}

depend() {
	need net
}

update() {
	check ${IPSEC_HOSTS_UPDATE} || return 1
	flush
	load ${IPSEC_HOSTS_UPDATE}
	if [ $? -eq 0 ]; then
		mv ${IPSEC_HOSTS_UPDATE} ${IPSEC_HOSTS_STABLE}
	else
		eerror "Update failed - restarting"
		restart
	fi
}

start() {
	check ${IPSEC_HOSTS_STABLE} || return 1
	flush
	load ${IPSEC_HOSTS_STABLE}
	if [ $? -ne 0 ]; then
		eerror "Start failed - stopping"
		flush
		return 2
	fi
}

stop() {
	flush
}

check() {
	local IPSEC_HOSTS_FILE=$1
	local STATUS_CHECK=0
	ebegin "Checking ipsec hosts file ${IPSEC_HOSTS_FILE}"
	if [ -f "${IPSEC_HOSTS_FILE}" ]; then
		ERROR_SPI=`awk '/^[^#].*$/{ print $3 }' ${IPSEC_HOSTS_FILE} | sort | uniq -d`
		if [ -n "$ERROR_SPI" ];then
			eerror "SPI is not unique: ${ERROR_SPI}"
			STATUS_CHECK=3
		fi
		ERROR_IP=`awk '/^[^#].*$/{ print $1 }' ${IPSEC_HOSTS_FILE} | sort | uniq -d`
		if [ -n "$ERROR_IP" ];then
			eerror "IP is not unique: ${ERROR_IP}"
			STATUS_CHECK=4
		fi
	else
		eerror "ipsechosts file ${IPSEC_HOSTS_FILE} doesn't exist"
		STATUS_CHECK=2
	fi
	eend $STATUS_CHECK
}

isipv4() {
	local IPV4="$(echo $1 | awk '/[0-9]+\.[0-9]+\.[0-9]+\.[0-9]+/ { print $1 }')"

	if [ -z "$IPV4" ]
	then
		return 1
	else
		return 0
	fi
}

load() {
	local ipsechosts="${1}"

	local inets="$(ip a | sed -n -e 's/^.*\(inet\|inet6\) \([0-9a-f:.]*\)\/.*$/;\2;/p')"
	local hosts="$(awk '$1 ~ /^[0-9a-f:.]+$/{ printf ";%s;%04x%s\n", $1, strtonum("0x" $3), $4 }' "${ipsechosts}")"

	local  lefts="$(echo "${hosts}" | grep -F    "${inets}" | sed -e 's|^;||')"
	local rights="$(echo "${hosts}" | grep -F -v "${lefts}" | sed -e 's|^;||')"

	for _dir in 'in' 'out'; do
		for _desc in 'proto ipv6-icmp' 'dport 22' 'sport 22'; do
			ip xfrm policy add dir "${_dir}" src '::/0'      dst '::/0'      ${_desc} priority 2147483548
		done
		for _desc in 'proto icmp' 'dport 22' 'sport 22'; do
			ip xfrm policy add dir "${_dir}" src '0.0.0.0/0' dst '0.0.0.0/0' ${_desc} priority 2147483548
		done
	done

	for left in ${lefts}; do
		local leftip6=false
		isipv4 ${left%%;*} || leftip6=true
		for right in ${rights}; do
			local rightip6=false
			isipv4 ${right%%;*} || rightip6=true
			if [ $leftip6 == $rightip6 ]; then
				setxfrm "${left}" "${right}"
			fi
		done
	done
}

setxfrm() {
	local lefthost="${1%%;*}"
	local leftkey="${1##*;}"
	local righthost="${2%%;*}"
	local rightkey="${2##*;}"
	ip xfrm state add \
		src "${lefthost}" dst "${righthost}" \
		proto esp spi "0x${leftkey:0:4}${rightkey:0:4}" \
		mode transport aead "rfc4106(gcm(aes))" "0x${rightkey:4}" 96
	ip xfrm state add \
		src "${righthost}" dst "${lefthost}" \
		proto esp spi "0x${rightkey:0:4}${leftkey:0:4}" \
		mode transport aead "rfc4106(gcm(aes))" "0x${leftkey:4}" 96
	ip xfrm policy add dir out \
		src "${lefthost}" dst "${righthost}" \
		priority 2147483648 tmpl proto esp mode transport
	ip xfrm policy add dir in \
		src "${righthost}" dst "${lefthost}" \
		priority 2147483648 tmpl proto esp mode transport
}

flush() {
	ebegin "Flushing associations and policies"
	ip xfrm state flush
	ip xfrm policy flush
	eend $?
}
